# system/btrfs-scrub@.service
[Unit]
Description=BTRFS scrub %f

[Service]
ExecStart=/usr/bin/btrfs scrub start -Bd %f
SyslogIdentifier=btrfs-scrub@%f
