# system/pacman-Fy.service
[Unit]
Description=Sync pacman(8) db to allow local search (pacman -F)
Documentation=man:pacman(8)
After=reachable-retry@archlinux.org
Requires=reachable-retry@archlinux.org

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman -Fy
