# system/sbsign-kernel@.path
[Unit]
Description=Watch for changes on the %i kernel file to sign it

[Path]
PathChanged=/boot/%i

[Install]
WantedBy=multi-user.target
