# user/mbsync-a.service
[Unit]
Desctiption=Sync all mailboxes
Documentation=man:mbsync(1)

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -a
