# user/sctw.service
[Unit]
Description=screen color temperature wrapper

[Service]
Type=oneshot
ExecStart=/usr/bin/sctw
